# configs/server/systemd-user/openclaw-gateway.service
[Unit]
Description=OpenClaw Gateway (v2026.3.7)
After=network-online.target
Wants=network-online.target

[Service]
ExecStart=/usr/bin/node /home/tup/.npm-global/lib/node_modules/openclaw/dist/index.js gateway --port 18789
Restart=always
RestartSec=5
TimeoutStopSec=30
TimeoutStartSec=30
SuccessExitStatus=0 143
KillMode=control-group
Environment=HOME=/home/tup
Environment=TMPDIR=/tmp
Environment=PATH=/home/tup/.local/bin:/home/tup/.npm-global/bin:/home/tup/bin:/home/tup/.volta/bin:/home/tup/.asdf/shims:/home/tup/.bun/bin:/home/tup/.nvm/current/bin:/home/tup/.fnm/current/bin:/home/tup/.local/share/pnpm:/usr/local/bin:/usr/bin:/bin
Environment=OPENCLAW_GATEWAY_PORT=18789
Environment=OPENCLAW_SYSTEMD_UNIT=openclaw-gateway.service
Environment="OPENCLAW_WINDOWS_TASK_NAME=OpenClaw Gateway"
Environment=OPENCLAW_SERVICE_MARKER=openclaw
Environment=OPENCLAW_SERVICE_KIND=gateway
Environment=OPENCLAW_SERVICE_VERSION=2026.3.7

[Install]
WantedBy=default.target
